# deploy/researchpipe-frontend.service
[Unit]
Description=ResearchPipe Frontend (Next.js)
After=network.target

[Service]
Type=simple
User=muye
WorkingDirectory=/home/muye/researchpipe/frontend
Environment="NODE_ENV=production"
Environment="PORT=3726"
Environment="HOSTNAME=127.0.0.1"
EnvironmentFile=/home/muye/researchpipe/frontend/.env.local
ExecStart=/usr/bin/npm run start
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
