# deploy/researchpipe-backend.service
[Unit]
Description=ResearchPipe API (FastAPI)
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=muye
WorkingDirectory=/home/muye/researchpipe/backend
Environment="PYTHONPATH=/home/muye/researchpipe/backend/src"
EnvironmentFile=/home/muye/researchpipe/backend/.env
ExecStart=/home/muye/researchpipe/backend/.venv/bin/uvicorn researchpipe_api.main:app \
    --host 127.0.0.1 --port 3725 \
    --workers 2 \
    --timeout-keep-alive 60 \
    --proxy-headers --forwarded-allow-ips=127.0.0.1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
